# runtime/conf.d/dbus.conf
[program:dbus]
autorestart=true
command=/usr/bin/dbus-daemon --system --nofork
priority=5
stdout_logfile=/var/log/dbus.log
stderr_logfile=/var/log/dbus.err
