# runtime/conf.d/x11vnc.conf
[program:x11vnc]
autorestart=true
command=x11vnc -display :1 -forever -shared -xkb -repeat -capslock -rfbport 5900
priority=20
stdout_logfile=/var/log/x11vnc.log
stderr_logfile=/var/log/x11vnc.err
